# crontab
* 0,6,12,18 * * * /photoprism/cron/index.sh 2>> /var/log/photoprism-cron.log
